# scroff-server/deploy/scroff-server.service
[Unit]
Description=Scroff Server - ADB-based screen controller (JDK 21)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=scroff
WorkingDirectory=/opt/scroff-server
# 使用 local profile，自动加载同目录的 application-local.yml
# 敏感信息（远程 DB 地址/账号、ADB 路径）都写在 application-local.yml 里
# 显式指向 JDK 21 的 java，避免系统装了多个 JDK 时跑错版本
ExecStart=/usr/lib/jvm/java-21-openjdk-amd64/bin/java -jar /opt/scroff-server/scroff-server.jar --spring.profiles.active=local
Restart=always
RestartSec=5
StandardOutput=append:/var/log/scroff-server.log
StandardError=append:/var/log/scroff-server.err.log
Environment=JAVA_OPTS=-Xms256m -Xmx512m

# 安全加固（测试稳定后建议开启）
# NoNewPrivileges=true
# ProtectSystem=strict
# ReadWritePaths=/opt/scroff-server /var/log

[Install]
WantedBy=multi-user.target
